# cms/config/.lxc/system.d/craft-email-custom-entries-export-as-csv.timer
[Unit]
Description=Timer for scheduling service

[Timer]
OnCalendar=Thu 07:00
Persistent=true

[Install]
WantedBy=timers.target